# kde/autostart/ssh-add.desktop
[Desktop Entry]
Type=Application
Name=ssh-add
Comment=Load default SSH keys into the agent at KDE login (ksshaskpass prompts, KWallet stores the passphrase)
Exec=ssh-add
X-KDE-AutostartScript=true
OnlyShowIn=KDE;
Terminal=false
